# data/latecall.service
[Unit]
Description=Latecall service

[Service]
ExecStart=/usr/local/bin/latecall --revive
Restart=always

[Install]
WantedBy=multi-user.target